# config/gvmd.service.in
[Unit]
Description=Greenbone Vulnerability Manager daemon (gvmd)
After=network.target networking.service postgresql.service ospd-openvas.service
Wants=postgresql.service ospd-openvas.service
Documentation=man:gvmd(8)
ConditionKernelCommandLine=!recovery

[Service]
Type=forking
PIDFile=/usr/var/run/gvmd.pid
ExecStart=/usr/sbin/gvmd
KillMode=mixed
Restart=always
RestartSec=1
User=root
Group=root
TimeoutStartSec=600
TimeoutStopSec=60

[Install]
WantedBy=multi-user.target
